# dot_config/systemd/user/dms.service.d/override.conf
[Service]
ExecStart=
ExecStart=%h/.local/bin/dms-run-session
